# nginx-reload.timer
[Unit]
Description=Reload nginx daily to account for wildcard cert updates

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
