# nginx-reload.service
[Unit]
Description=Reload nginx service

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user reload nginx.service

