# pageboard-reload.service
[Unit]
Description=Restart pageboard service

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart pageboard.service

