# pageboard-reload.timer
[Unit]
Description=Reload leaking pageboard daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
